# archiso-troya/airootfs/usr/lib/systemd/system/startup-script.service
[Unit]
Description=Startup script

[Service]
Type=oneshot
ExecStart=/usr/bin/startup.sh
After=network-online.target
Wants=network-online.target
